# files/default/nix-daemon.conf
# vi:ft=upstart
description "Nix Daemon"
start on filesystem
stop on shutdown
respawn
exec /nix/var/nix/profiles/default/bin/nix-daemon --daemon
